# deploy/signi-cms.service
# systemd unit for the signi-web CMS (Fastify: blog SSR, API, form intake).
# Install to /etc/systemd/system/signi-cms.service on the VPS, then:
#   systemctl daemon-reload && systemctl enable --now signi-cms

[Unit]
Description=Signi CMS (blog SSR + API)
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/srv/signi-web-cms
ExecStart=/usr/bin/node src/index.js
Restart=on-failure
RestartSec=3

# Staging configuration (dev.signi.com/prototypes/signi-web/).
# For production set BASE_PATH= (empty) and SITE_URL=https://signi.com.
Environment=NODE_ENV=production
Environment=HOST=127.0.0.1
Environment=PORT=8788
Environment=BASE_PATH=/prototypes/signi-web
Environment=SITE_URL=https://dev.signi.com
Environment=DB_PATH=/var/lib/signi-web/signi.sqlite
Environment=MAIL_FROM=Signi <noreply@signi.com>
Environment=SALES_INBOX=sales@signi.com
# Secrets — set when activated. Without RESEND_API_KEY form mail is logged only.
# Environment=RESEND_API_KEY=
# Environment=TURNSTILE_SECRET=

[Install]
WantedBy=multi-user.target
